# docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/data/logs

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:web]
command=/usr/local/bin/python /app/app.py
directory=/app
autostart=true
autorestart=true
startsecs=2
priority=10
stdout_logfile=/data/logs/web-supervisor.log
stderr_logfile=/data/logs/web-supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[program:cleaner]
command=/app/docker/run_cleaner.sh
directory=/app
autostart=true
autorestart=true
startsecs=2
priority=20
stdout_logfile=/data/logs/cleaner-supervisor.log
stderr_logfile=/data/logs/cleaner-supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
